# debian/prerm
#!/bin/bash
# postrm script

set -e

case "$1" in
    remove|purge)
        if command -v supervisorctl >/dev/null 2>&1; then
            supervisorctl reread
            supervisorctl update
        else
            systemctl stop nmexec.service || true
            systemctl disable nmexec.service || true
            systemctl daemon-reexec || true
        fi

        userdel nmexec || true
        groupdel nmexec || true

        rm -rf /var/lib/nmexec
    ;;

    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0